# config/systemd/user/wallpaper.timer
[Unit]
Description=Randomize wallpapers every 40 minutes

[Timer]
OnBootSec=0min
OnUnitActiveSec=40min

[Install]
WantedBy=default.target
